# install/launch_flask_at_startup.conf
[program:farm_watcher_flask]
command=/home/pi/farm_watcher/launch_flask.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/farm_watcher.err.log
stdout_logfile=/var/log/farm_watcher.out.log

